# the-print-farm.service
[Unit]
Description=The Print Farm Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/the-print-farm
ExecStart=/root/the-print-farm/venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
